# templates/systemd/etcd.service.j2
{#
   Copyright 2015 VMware, Inc.

   Licensed under the Apache License, Version 2.0 (the "License");
   you may not use this file except in compliance with the License.
   You may obtain a copy of the License at

       http://www.apache.org/licenses/LICENSE-2.0

   Unless required by applicable law or agreed to in writing, software
   distributed under the License is distributed on an "AS IS" BASIS,
   WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
   See the License for the specific language governing permissions and
   limitations under the License.
#} 
[Unit]
Description=etcd

[Service]
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/etcd --advertise-client-urls http://{{ master_address }}:2379 --listen-client-urls http://0.0.0.0:2379 --initial-advertise-peer-urls http://{{ master_address }}:2380  --listen-peer-urls http://0.0.0.0:2380  --initial-cluster-token etcd-cluster-1  -initial-cluster default=http://{{ master_address }}:2380  -initial-cluster-state new
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
